# govtchat.service
# /etc/systemd/system/govtchat.service
#
# Systemd unit file for the Chat Agent API Service
# Corrected for user 'vpa' with Conda environment 'cogops'.

[Unit]
Description=Chat Agent API Service
# Ensures the network is available before attempting to start.
After=network.target

[Service]
# --- User and Environment Configuration ---
# Run the service as the 'vpa' user.
User=vpa
Group=vpa

# The absolute path to your project's root directory.
# This is crucial for relative paths in your code to work correctly.
WorkingDirectory=/home/vpa/CogOpsCB

# The absolute path to your project's .env file.
# This will load all your environment variables (API keys, etc.).
EnvironmentFile=/home/vpa/CogOpsCB/.env

# --- Execution Command ---
# The command to start the FastAPI application.
# It uses the absolute path to the uvicorn executable within your specific Conda environment.
# We also specify the file name 'api:app' which is standard.
# The --workers flag runs multiple processes to handle more concurrent requests.
ExecStart=/home/vpa/miniconda3/envs/cogops/bin/uvicorn api:app --host 0.0.0.0 --port 9000 --workers 4

# --- Service Behavior ---
# Always restart the service if it stops unexpectedly (e.g., crashes).
Restart=always
# Wait 10 seconds before attempting to restart.
RestartSec=10

# --- Logging ---
# Redirect stdout and stderr to the system's journal.
StandardOutput=journal
StandardError=journal
# A unique identifier for this service's logs in the system journal.
SyslogIdentifier=chat-agent-api

[Install]
# Enable the service to start automatically on system boot.
WantedBy=multi-user.target