# evaluation/retriver-evaluation.service
# /etc/systemd/system/retrieval-evaluation.service
#
# Systemd unit file for the Retrieval Evaluation API Service.
# Configured for user 'vpa' with Conda environment 'cogops'.

[Unit]
Description=Retrieval Evaluation API Service
After=network.target

[Service]
# --- User and Environment Configuration ---
User=vpa
Group=vpa

# The absolute path to your project's root directory.
# This is crucial for Python's import system to find your modules.
WorkingDirectory=/home/vpa/CogOpsCB

# The absolute path to your project's .env file to load all secrets.
EnvironmentFile=/home/vpa/CogOpsCB/.env

# --- Execution Command ---
# Runs the uvicorn server from the specific conda environment.
# 'evaluation.api_server:app' is the Python module path to your FastAPI app object.
# --workers 4 is a good starting point for handling concurrent requests in production.
ExecStart=/home/vpa/miniconda3/envs/cogops/bin/uvicorn evaluation.api_server:app --host 0.0.0.0 --port 9051

# --- Service Behavior ---
# Always restart the service if it stops unexpectedly.
Restart=always
RestartSec=10

# --- Logging ---
# Redirect all output to the system's journal for easy inspection.
StandardOutput=journal
StandardError=journal
# A unique identifier for this service's logs.
SyslogIdentifier=retrieval-eval-api

[Install]
# This allows the service to start automatically on system boot.
WantedBy=multi-user.target