# search_and_crawl/crawler.service
[Unit]
Description=Intelligent Crawler API Service
After=network.target

[Service]
User=vpa
Group=vpa

# The directory where your project is located
WorkingDirectory=/home/vpa/CogOpsCB/search_and_crawl

# The file containing your secrets (SEARXNG_API_URL, SEARXNG_API_KEY)
EnvironmentFile=/home/vpa/CogOpsCB/search_and_crawl/.env

# The command to start the Uvicorn server
# We use the full, absolute path to the uvicorn executable
ExecStart=/home/vpa/miniconda3/bin/uvicorn main:app --host 0.0.0.0 --port 9234

# Restart the service automatically if it fails
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target